# deploy/balupi.service
[Unit]
Description=BaluPi - Smart Cache & Energy Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/balupi/backend
ExecStart=/opt/balupi/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Environment
EnvironmentFile=/opt/balupi/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/balupi/data /opt/balupi/backend/data
PrivateTmp=true

# Watchdog
WatchdogSec=30

[Install]
WantedBy=multi-user.target
